# supervisor.service
[Unit]
Description=Supervisor Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /absolute/path/to/your/project/bin/supervisor
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=supervisor
User=supervisor
Group=supervisor

[Install]
WantedBy=multi-user.target
