# supervisor.service
[Unit]
Description=Supervisor Service
After=network.target

[Service]
ExecStart=/bin/sh -c '/usr/bin/python3 /usr/local/bin/supervisor 2>&1 | /usr/bin/tee -a /var/log/thirdreality/supervisor.log'
Restart=on-failure
RestartSec=5
# Prepare log directory and rotate previous log to .1 on each start
ExecStartPre=/bin/sh -c 'mkdir -p /var/log/thirdreality'
ExecStartPre=/bin/sh -c '[ -f /var/log/thirdreality/supervisor.log ] && mv /var/log/thirdreality/supervisor.log /var/log/thirdreality/supervisor.log.1 || true'
ExecStartPre=/bin/sh -c 'touch /var/log/thirdreality/supervisor.log && chmod 0644 /var/log/thirdreality/supervisor.log'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=supervisor
# Removed User and Group settings to run as root
# User=supervisor
# Group=supervisor

[Install]
WantedBy=multi-user.target
